# rootfs/etc/services.d/notifier/run
#!/usr/bin/with-contenv bash

while :
do
	/usr/bin/curl --connect-timeout 10 $NOTIFY_URL
	sleep $NOTIFY_EVERY_SECONDS
done